# statlang/crontab.cron
00 01 * * * python3.6 /home/ununtu/predictTrend/Statistic-Lang/statlang/dbug.py
